# Djangodemo/conf/supervisor_celery_worker.ini
[program:celery-worker]
directory=/home/jerd/PycharmProjects/Djangodemo
command=python manage.py celery worker -l INFO
environment=PATH='/home/jerd/env/envpy3/bin'
stdout_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.worker.log
stderr_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.worker.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998