# Djangodemo/conf/supervisor_celery_flower.ini
[program:celery-flower]
directory=/home/jerd/PycharmProjects/Djangodemo
command=python manage.py celery flower --address=0.0.0.0 --port=5555
environment=PATH='/home/jerd/env/envpy3/bin'
stdout_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.flower.log
stderr_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=995