# Djangodemo/conf/supervisor_celery_beat.ini
[program:celery-beat]
directory=/home/jerd/PycharmProjects/Djangodemo
command=python manage.py celery beat -l INFO
environment=PATH='/home/jerd/env/envpy3/bin'
stdout_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.beat.log
stderr_logfile=/home/jerd/PycharmProjects/Djangodemo/logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=997